# ops/services/biketag-db-backup.service
[Unit]
Description=Collision DB Backup

[Service]
TimeoutStartSec=0
ExecStart=/home/core/collision/ops/bin/archive-db 

