# ops/services/collision-db-backup.timer
[Unit]
Description=Run Collision DB Backups Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
